# configs/systemd/hsaj-core.service
[Unit]
Description=HSAJ core operator API and runtime jobs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/hsaj/hsaj.env
Environment=HSAJ_ROOT=/opt/hsaj
ExecStart=/usr/bin/env bash -c 'cd "${HSAJ_ROOT}/core" && exec hsaj serve'
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
