# configs/systemd/hsaj-core.timer
[Unit]
Description=HSAJ legacy scheduled smoke (scan + roon sync + apply --dry-run)

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=15m
Persistent=true
Unit=hsaj-maintenance.service

[Install]
WantedBy=timers.target
